# ops/config/cerberus-sdr.service
[Unit]
Description=Nexa Cerberus SDR agent (always-on)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cerberus
EnvironmentFile=-/opt/configs/env
ExecStart=/opt/cerberus/cerberus agent --config /opt/configs/sdr-agent.json
Restart=always
RestartSec=10
StandardOutput=append:/var/log/cerberus/sdr.log
StandardError=append:/var/log/cerberus/sdr-error.log

[Install]
WantedBy=multi-user.target
