# ops/config/nexa-web.service
[Unit]
Description=Nexa Web Interface
After=network.target

[Service]
Type=simple
User=root
Group=root
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3003
EnvironmentFile=-/opt/aura/apps/web/shared/.env.production
WorkingDirectory=/opt/aura/apps/web/current
ExecStart=/usr/bin/env node .next/standalone/server.js
Restart=always
RestartSec=2
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
