# ops/config/nexa_autopilot.service
[Unit]
Description=Nexa Autopilot Control Loop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root
Environment=AURA_ROOT=/root
Environment=NEXA_GATEWAY_PORT=9080
Environment=NEXA_DOCS_INTERVAL_SEC=1800
Environment=NEXA_OPS_CAST_INTERVAL_SEC=3600
Environment=NEXA_OPS_CAST_SOURCE=agency_metrics
Environment=NEXA_OPS_CAST_PROVIDER=local
ExecStart=/usr/bin/python3 /root/ops/autopilot/nexa_autopilot.py loop --interval 60
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
