# ops/config/nexa-gateway.service
[Unit]
Description=Nexa Gateway
After=network.target

[Service]
Type=simple
User=root
EnvironmentFile=-/etc/nexa/runtime.env
ExecStart=/usr/bin/bash -lc 'export NEXA_ROOT="${NEXA_ROOT:-/opt/nexa/runtime}"; export NEXA_REPO_DIR="${NEXA_REPO_DIR:-/opt/nexa/repo}"; export AURA_GATEWAY_PORT="${AURA_GATEWAY_PORT:-${NEXA_GATEWAY_PORT:-8765}}"; cd "$NEXA_REPO_DIR/core/nexa-gateway"; exec ./zig-out/bin/nexa-gateway'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
