# llm-port-project/src/aura/sovereign-stack/aura.service
[Unit]
Description=Aura Mesh Node
[Service]
ExecStart=/home/yani/Aura/bin/aura-mesh daemon
