# apps/wealth/aura-trader.service
[Unit]
Description=Aura Algo Trader Daemon
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
User=yani
WorkingDirectory=/home/yani/Aura/ai_agency_wealth
ExecStart=/home/yani/Aura/ai_agency_wealth/venv312/bin/python3 trader_daemon.py
ExecStop=/home/yani/Aura/ai_agency_wealth/venv312/bin/python3 trader_ctl.py stop

# Restart on crash, but not on clean stop (exit code 0)
Restart=on-failure
RestartSec=10

# Environment
EnvironmentFile=/home/yani/Aura/ai_agency_wealth/.env

# Logging (goes to journald; also rotated file in logs/)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aura-trader

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ReadWritePaths=/home/yani/Aura/ai_agency_wealth

[Install]
WantedBy=multi-user.target
